# ansible/roles/hosts/files/etc/systemd/system/ansible.timer
[Unit]
Description=Call ansible

[Timer]
OnCalendar=*-*-* 00:00:00
Unit=ansible.service

[Install]
WantedBy=timers.target
