# ansible/roles/hosts/files/etc/systemd/system/ansible.service
[Unit]
Description=Call myos ansible-pull

[Service]
Type=oneshot
ExecStart=/etc/init.d/myos ansible-pull

[Install]
WantedBy=multi-user.target
