# ansible/roles/hosts/files/etc/periodic/hourly/ansible
#!/bin/sh

/etc/init.d/myos ansible-pull > /var/log/ansible.log
