# ansible/roles/hosts/files/etc/init.d/myos
#!/bin/bash
# shellcheck source=/dev/null disable=SC2046
### BEGIN INIT INFO
# Provides: myos
# Required-Start: docker
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: MYOS system calls
### END INIT INFO

NAME="$(basename "$0")"

# check system config
[ -r /etc/default/myos ] && . /etc/default/myos 2>/dev/null
type make >/dev/null 2>&1 && [ -d "${MYOS}" ] && cd "${MYOS}" || exit 1

case "$1" in
  force-stop)
    HOOK="stop stack-node-stop"
    shift
    ;;
  force-reload)
    HOOK="recreate stack-node-recreate"
    shift
    ;;
  force-restart)
    HOOK="restart stack-node-restart"
    shift
    ;;
  force-shutdown)
    HOOK="down stack-node-down"
    shift
    ;;
  reload)
    HOOK="recreate"
    shift
    ;;
  shutdown)
    HOOK="down"
    shift
    ;;
  start)
    echo "Starting $NAME..."
    HOOK="up"
    shift
    ;;
  status)
    HOOK="ps"
    shift
    ;;
  stop)
    echo "Stopping $NAME..."
    ;;
  *)
    ;;
esac

IFS=$'\n'; exec env $(cat /etc/default/myos ~/.myos 2>/dev/null) make ${HOOK:-} "${@}"
