# ansible/roles/hosts/files/etc/init.d/myos
#!/bin/bash
# shellcheck source=/dev/null disable=SC2046
### BEGIN INIT INFO
# Provides: myos
# Required-Start: docker
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: MYOS system calls
### END INIT INFO

NAME="$(basename "$0")"

# check system config
[ -r /etc/default/myos ] && . /etc/default/myos 2>/dev/null
type make >/dev/null 2>&1 && [ -d "${MYOS}" ] && cd "${MYOS}" || exit 1

case "$1" in
  force-reload)
    TARGET='force-recreate'
    ;;
  reload)
    TARGET="recreate"
    ;;
  start)
    echo "Starting $NAME..."
    TARGET="up"
    ;;
  status)
    TARGET="ps"
    ;;
  stop)
    echo "Stopping $NAME..."
    TARGET="$1"
    ;;
  *)
    TARGET="$1"
    ;;
esac
shift ||:

IFS=$'\n'; exec env $(cat /etc/default/myos 2>/dev/null) make ${TARGET:-} ${TARGET_FORCE:-} $@
