# ansible/roles/hosts/files/etc/cron.d/ansible
PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin

@reboot   root /etc/init.d/myos ansible-pull > /var/log/ansible.log
0 * * * * root /etc/init.d/myos ansible-pull > /var/log/ansible.log
